# thanos/thanos-sidecar.service
# -*- mode: conf -*-

[Unit]
Description=Thanos Sidecar service.
Documentation=https://thanos.io/tip/thanos/
After=network.target
After=prometheus.service

[Service]
EnvironmentFile=-/etc/default/thanos-sidecar
User=prometheus
ExecStart=/usr/bin/thanos sidecar \
          $THANOS_SIDECAR_OPTS
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
